# challenges/web/medium/monitoring-dashboard/src/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:icinga2]
command=/usr/sbin/icinga2 daemon
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=nagios
redirect_stderr=false

[program:proxy]
command=/usr/local/bin/proxy
directory=/app
environment=ROCKET_ADDRESS="0.0.0.0",ROCKET_PORT="8000"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
redirect_stderr=false
